# systemd/user/docker.target
[Unit]
Description=Docker Target
# StopWhenUnneeded=yes